# Linux/Packages/RPM-build/php-fpm
#! /bin/sh
# Comments to support chkconfig on CentOS
# chkconfig: 2345 65 37
#
set -e

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DESC="php-fpm daemon"
NAME=php-fpm

#php-fpm命令的路径
DAEMON=/usr/sbin/$NAME

#php-fpm的配置文件
CONFIGFILE=/usr/local/php-7.1.11/etc/php-fpm.conf

#php-fpm的pid文件
PIDFILE=/var/run/php-fpm/$NAME.pid

SCRIPTNAME=/etc/init.d/$NAME

# Gracefully exit if the package has been removed.
test -x $DAEMON || exit 0
d_start() { 
    if netstat -ntl|grep 9000 &>/dev/null;then
      echo  " already running"
    else
       $DAEMON -y $CONFIGFILE
       pid_id=`ps aux | grep "$CONFIGFILE" |grep -v 'grep' |awk '{print $2}'`
       echo $pid_id > $PIDFILE
       echo -n "Starting $DESC is success" 
    fi
}

d_stop() {
     if netstat -ntl|grep 9000 &>/dev/null;then
       p_id=`cat $PIDFILE`
       kill $p_id
       rm -f $PIDFILE
       echo -n "Stopping $DESC is success" 
     else
       echo "not running"
     fi
}

d_reload() { 
   d_stop
   d_start
   echo -n "Restarting $DESC is success" 
}

case "$1" in 
  start) 
          d_start 
          echo "." ;; 
  stop) 
          d_stop 
          echo "." ;; 
  reload) 
          echo -n "Reloading $DESC configuration..." 
          d_reload 
          echo "reloaded." ;; 
    restart) 
          d_stop 
          sleep 1 
          echo 
          d_start 
          echo "." ;; 
    *) 
          echo "Usage: $SCRIPTNAME {start|stop|restart|force-reload}" >&2 
          exit 3 ;;
esac


